# dashboard_v3/ops/systemd/trading-dashboard-v3.service
[Unit]
Description=Trading Dashboard V3 Sidecar
After=network.target

[Service]
Type=simple
User=__SERVICE_USER__
WorkingDirectory=__BOT_ROOT__/dashboard_v3/backend
ExecStart=__BOT_ROOT__/dashboard_v3/backend/.venv/bin/python __BOT_ROOT__/dashboard_v3/backend/run.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
