# dashboard_v3/ops/systemd/trading-dashboard-v3-tail.service
[Unit]
Description=Trading Dashboard V3 Bot Log Tail
After=docker.service trading-dashboard-v3.service
Requires=trading-dashboard-v3.service

[Service]
Type=simple
User=root
WorkingDirectory=__BOT_ROOT__
ExecStart=/bin/bash -lc 'LOG_FILE=$(docker inspect -f "{{.LogPath}}" __BOT_CONTAINER__); exec __BOT_ROOT__/dashboard_v3/backend/.venv/bin/python __BOT_ROOT__/dashboard_v3/agents/log_tail_agent.py --log "$LOG_FILE" --dashboard http://127.0.0.1:8000 --from-start'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
