# ops/systemd/trading-dashboard.service
[Unit]
Description=Trading Portfolio Visual Dashboard
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/opt/trading-bot/dashboard/backend
ExecStart=/opt/trading-bot/dashboard/backend/.venv/bin/python /opt/trading-bot/dashboard/backend/run.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
