# dashboard_v2/ops/systemd/trading-dashboard-v2-tail.service
[Unit]
Description=Trading Dashboard V2 Docker Log Tail
After=trading-dashboard-v2.service docker.service
Requires=trading-dashboard-v2.service

[Service]
Type=simple
User=root
WorkingDirectory=__BOT_ROOT__
Environment=BOT_CONTAINER=__BOT_CONTAINER__
ExecStart=/bin/bash -lc 'LOG_FILE=$(docker inspect -f "{{.LogPath}}" "$BOT_CONTAINER"); exec __BOT_ROOT__/dashboard_v2/backend/.venv/bin/python __BOT_ROOT__/dashboard_v2/agents/log_tail_agent.py --log "$LOG_FILE" --dashboard http://127.0.0.1:8000 --from-start'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
